# worker-scheduler.conf
[supervisord]
nodaemon=true

[program:scheduler]
command=/app/bin/docker-entrypoint scheduler

[program:worker]
command=/app/bin/docker-entrypoint worker